# imageroot/systemd/controller-metrics-exporter.path
[Unit]
Description=Watch for vpn connections from MODULE_ID-vpn.service
BindsTo=MODULE_ID.service
After=MODULE_ID-vpn.service

[Path]
PathChanged=/var/lib/nethserver/MODULE_ID/state/clients
