# imageroot/systemd/user/metrics-exporter.path
[Unit]
Description=Watch for vpn connections from vpn.service
BindsTo=controller.service
After=vpn.service

[Path]
PathChanged=%S/state/clients
